# runcat/systemd/runcat-waybar.service
[Unit]
Description=RunCat helper for Waybar
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/runcat/scripts/runcat.py --interval 0.4
Restart=on-failure

[Install]
WantedBy=default.target
